# etc/supervisor/conf.d/video_looper.conf
# Supervisord configuration to run video looper at boot and
# ensure it runs continuously.
[program:video_looper]
command=python3 -u -m Adafruit_Video_Looper.video_looper
stdout_logfile=/var/log/MF/mf.log
stderr_logfile=/var/log/MF/mfError.log
autostart=true
autorestart=unexpected
startsecs=5

[program:magikflame_controls]
command=/usr/bin/python3 /home/pi/server.py
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/var/log/MF/server.log
stderr_logfile=/var/log/MF/serverError.log

